# supervisor.conf
[program:logstash]
user=root
command=/scripts/run.sh
redirect_stderr=true
autorestart=true
